# wolfgang/com.crespo.wolfgang.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.crespo.wolfgang</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/crespo/wolfgang/venv/bin/python</string>
        <string>/Users/crespo/wolfgang/wolfgang_server.py</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/crespo/wolfgang/wolfgang.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/crespo/wolfgang/wolfgang.err</string>
    <key>WorkingDirectory</key>
    <string>/Users/crespo/wolfgang</string>
</dict>
</plist>
